# deploy/qadam-backend.service
[Unit]
Description=Qadam Backend (Phoenix)
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/qadam/qadam_backend
EnvironmentFile=/home/ubuntu/qadam/.env.production
Environment=PATH=/home/ubuntu/.asdf/shims:/home/ubuntu/.asdf/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/ubuntu/.asdf/shims/mix phx.server
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qadam-backend

[Install]
WantedBy=multi-user.target
